# pi-radio.service
# /etc/systemd/system/pi-radio.service
# 自動起動設定:
#   sudo cp pi-radio.service /etc/systemd/system/
#   sudo systemctl enable pi-radio
#   sudo systemctl start pi-radio

[Unit]
Description=Pi Radio - Internet Radio Server
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/netradio
ExecStart=/usr/bin/python3 /home/pi/netradio/app.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
# ALSA 出力先（必要に応じて変更）
# Environment=ALSA_CARD=0

[Install]
WantedBy=multi-user.target
